# scripts/cron.d/index_update.sh
#!/bin/bash                                                     
*/15 * * * * root /usr/bin/python /home/gestion/bomberos/manage.py update_index >> /var/log/syslog
